# ddns.service
[Unit]
Description=Compare current IP and update
Requisite=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/ddns
